# assets/wireplumber.service.d/reflash.conf
[Service]
ExecStartPost=/usr/local/bin/mark2-xvf-post-wp.sh
